# systemd/voice-transcribe-server.service
[Unit]
Description=Voice Transcription Server (Parakeet on GPU)
After=default.target

[Service]
Type=simple
ExecStart=/home/vijay-katta/faster-whisper-env/bin/python3 /home/vijay-katta/bin/voice-transcribe-server
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
Environment=NEMO_LOG_LEVEL=CRITICAL
Environment=HYDRA_FULL_ERROR=0
Environment=HOME=/home/vijay-katta
StandardOutput=append:/tmp/voice-transcribe-server.log
StandardError=append:/tmp/voice-transcribe-server.log

[Install]
WantedBy=default.target
